# root/etc/services.d/caddy/run
#!/usr/bin/with-contenv bash

# reference
# https://github.com/caddyserver/forwardproxy
# https://github.com/caddyserver/forwardproxy/blob/master/docker-build/gen_caddyfile_and_start.sh

if [ "${PROXY_ENABLED}" != "true" ]; then
    s6-svc -d .
    exit 0
fi

command="caddy -agree -conf /config/Caddyfile"
command="$command${PROXY_VERBOSE+ -log stdout}"

echo "*** executing => $command"
exec s6-setuidgid abc $command
