# root/etc/services.d/gt/run
#!/usr/bin/with-contenv bash

export HOME=/config

if [ "$(printenv GT_USE | tr -d '"' | tr -d "'")" == "true" ]; then

    command="gt --ip 0.0.0.0 --port ${GT_PORT:-21000} --system-proxy false --silent"
    command="$command${GT_VERBOSE+ --verbose \"green-tunnel:*\"}"
    command="$command --dns-type ${GT_DNSTYPE:-https}"
    command="$command --dns-server ${GT_DNSSERVER:-https://1.1.1.1/dns-query}"

    echo "running => $command"
    exec s6-setuidgid abc $command

fi
