# root/etc/services.d/gt/run
#!/usr/bin/with-contenv bash

if [ "$(printenv USE_GT | tr -d '"' | tr -d "'")" == "true" ]; then

    command="gt --ip 127.0.0.1 --port 21000"
    command="$command --dnsType ${GT_DNSTYPE:-https}"
    command="$command --dnsServer ${GT_DNSSERVER:-https://cloudflare-dns.com/dns-query}"

    echo ""
    echo "running => $command"
    echo ""
    
    exec s6-setuidgid abc $command

fi
